# ops/systemd/stackchan-cloud-news.service
[Unit]
Description=StackChan Cloud News MCP
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/stackchan-codex-bridge
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=8788
EnvironmentFile=-/etc/stackchan-cloud-news.env
ExecStart=/usr/bin/npm run start:cloud-news
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
