# ops/systemd/stackchan-xiaozhi-watchdog.service
[Unit]
Description=StackChan Xiaozhi Watchdog
After=stackchan-cloud-news.service stackchan-xiaozhi-client.service

[Service]
Type=oneshot
ExecStart=/opt/stackchan-codex-bridge/ops/scripts/xiaozhi-watchdog.sh
